# rc.local
#!/bin/bash
openvt -s -w /opt/mylogon/logon.sh